# peach-dyndns-updater/debian/peach-dyndns-updater.timer
[Unit]
Description=Systemd timer which keeps a dynamic dns subdomain up to date with the latest device IP using nsupdate.

[Timer]
Unit=peach-dyndns-updater.service
OnCalendar=*:0/5
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
